# brightness-key/etc/systemd/system/hw_changes_brightness.service
[Unit]
Description=Set hw_changes_brightness

[Service]
ExecStart=/usr/local/bin/set_hw_changes_brightness.sh
